# sprinklers.cron
0 6 * * 2,4,6   root    "/usr/bin/python3 /var/sprinkler-controller/sprinkler-run.py -z 1 -d 15"
15 6 * * 2,4,6   root    "/usr/bin/python3 /var/sprinkler-controller/sprinkler-run.py -z 2 -d 15"
30 6 * * 2,4,6   root    "/usr/bin/python3 /var/sprinkler-controller/sprinkler-run.py -z 3 -d 15"
45 6 * * 2,4,6   root    "/usr/bin/python3 /var/sprinkler-controller/sprinkler-run.py -z 4 -d 15"
0 7 * * 2,4,6   root    "/usr/bin/python3 /var/sprinkler-controller/sprinkler-run.py -z 5 -d 15"
15 7 * * 2,4,6   root    "/usr/bin/python3 /var/sprinkler-controller/sprinkler-run.py -z 6 -d 15"
30 7 * * 2,4,6   root    "/usr/bin/python3 /var/sprinkler-controller/sprinkler-run.py -z 7 -d 15"
45 7 * * 2,4,6   root    "/usr/bin/python3 /var/sprinkler-controller/sprinkler-run.py -z 8 -d 15"